# system_admin/ansible_infrastructure/templates/pg_backup.cron.j2
# PostgreSQL backup schedule managed by Ansible
15 2 * * * postgres /usr/local/sbin/pg_backup.sh
